# systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -d 10 -c 20 -w 30 -f 100 -e
